# .config/sv/yt-local/run
#!/usr/bin/bash
#
#

# redirects all stderr output to stdout
exec 2>&1

source /home/amf/.local/opt/yt-local/venv/bin/activate; 
exec chpst -u amf python /home/amf/.local/opt/yt-local/server.py
