# .config/sv/pipewire/run
#!/usr/bin/bash
#

# redirects all stderr output to stdout
exec 2>&1

/usr/bin/pipewire &
/usr/bin/pipewire-pulse &

exec /usr/bin/wireplumber &
