# .config/sv/transmission-daemon/finish
#!/usr/bin/bash
#
#

exec echo "transmission-daemon finished"
