# .config/sv/transmission-daemon/run
#!/usr/bin/bash
#
#
#

exec chpst -u amf:transmission transmission-daemon -f --log-error
