# .config/sv/mpd/finish
#!/usr/bin/bash
#
#

exec echo "mpd finished"
