# .config/sv/mpd/run
#!/usr/bin/bash
#
#

# redirects all stderr output to stdout
exec 2>&1

exec chpst -u amf mpd --no-daemon
